# neptun-bridge/rootfs/etc/services.d/neptun-bridge/run
#!/usr/bin/with-contenv bash
set -e
exec /usr/local/bin/neptun_bridge.py
