# neptun-bridge/rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bash
exec mosquitto -c /etc/mosquitto/mosquitto.conf
